# rpi-meteo.service
[Unit]
Description=Meteo Service Dependencies
Requires=docker.service
After=docker.service

[Service]
#Restart=always
ExecStart=/usr/bin/docker-compose -f /home/pi/meteo/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /home/pi/meteo/docker-compose.yml down

[Install]
WantedBy=default.target
